# root/etc/services.d/tautulli/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid tautulli python3 /app/tautulli/Tautulli.py \
	--datadir /config
