# root/etc/services.d/tautulli/run
#!/usr/bin/with-contenv bash
uid=${PUID:-1000}
user=$(getent passwd $uid | awk -F: '{print $1}')

exec s6-setuidgid $user python3 /app/tautulli/Tautulli.py \
	--datadir /config
